# roles/common_fstrim/templates/systemd/fstrim-task.timer.j2
# **** {{ ansible_managed }} ****

[Unit]
Description=run fstrim daily

[Timer]
OnCalendar=*-*-* 4:15
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
