# roles/prometheus_node_exporter/templates/systemd/prometheus-node-generate-textfile-metrics-hourly.timer.j2
# **** {{ ansible_managed }} ****

[Unit]
Description=Generate textfile metrics for prometheus-node-exporter

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
