# roles/memcached/templates/systemd/memcached.target.j2
# **** {{ ansible_managed }} ****

[Unit]
Description=Memcached service group

[Install]
WantedBy=multi-user.target
