# roles/tinc/templates/systemd/tinc.target.j2
# **** {{ ansible_managed }} ****

[Unit]
Description=Tinc service group

[Install]
WantedBy=multi-user.target
