# roles/certbot/templates/systemd/certbot-task-renew.service.j2
# **** {{ ansible_managed }} ****

[Unit]
Description=Renew certificates

[Service]
Type=oneshot
ExecStart={{ certbot__cmd_new_base }}/certbot-renew
User=root
