# roles/certbot/templates/systemd/certbot-task-renew.timer.j2
# **** {{ ansible_managed }} ****

[Unit]
Description=Renew certificates

[Timer]
OnCalendar=Sun,Thu 3:15

[Install]
WantedBy=timers.target
